# packaging/debian/postinst
#!/bin/bash

systemctl enable proxyapp.service

exit 0
