# packaging/debian/prerm
#!/bin/bash

systemctl disable proxyapp.service

exit 0
