# proxyapp_py_ros2/packaging/debian/postinst
#!/bin/bash

systemctl enable proxyapp_py_unsecure_client.service

exit 0
