# proxyapp_py_ros2/packaging/debian/prerm
#!/bin/bash

systemctl disable proxyapp_py_secure_srv.service
systemctl disable proxyapp_py_unsecure_client.service

exit 0
